# supervisor/serv.conf
[program:webserver]
directory = /code
command = gunicorn project.wsgi:application --bind 0.0.0.0:8000

[program:huey]
directory = /code
command = python manage.py run_huey
